# 01-代码部分/1-01-现场代码/03_无锡远纺-1.1老版本/ARM1-1126-现场调整/rc.local
#!/bin/bash


sleep 12s

cd  /opt/app/ai-product-haijiang/aiInferServer

sudo HOME=/home/appuser nohup python3 manage.py runserver 0.0.0.0:8000 > /mnt/data/opt/logs/txl01/ai.log 2>&1 &

sleep 5s

mkdir -p /opt/logs/ai-product-haijiang/ai-cotton-api-server-release/

/opt/app/ai-product-haijiang/ai-cotton-api-server-release/autorun.sh &

sleep 60s

cd /opt/app/ai-product-haijiang/ai-hj-service/aiHjService/

#sudo HOME=/home/appuser nohup python3 ai_hj_run_service.py > /mnt/data/opt/logs/txl01/ctrl.log 2>&1 &
sudo HOME=/home/appuser nohup python3 daemon.py > /mnt/data/opt/logs/txl01/daemon.log 2>&1 &







